# deployments/ontolith-management-server.service
[Unit]
Description=Ontolith management server (L5 control plane)
Documentation=file:///home/ontolith/docs/L5-ontolith-access-security.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sharkyai
Group=sharkyai
WorkingDirectory=/home/ontolith

EnvironmentFile=-/etc/ontolith/ontolith-management.env
Environment=RUST_BACKTRACE=1

ExecStart=/usr/local/bin/ontolith-management-server
Restart=on-failure
RestartSec=3
TimeoutStopSec=20

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/ontolith /home/ontolith/data
AmbientCapabilities=
CapabilityBoundingSet=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ontolith-management-server

[Install]
WantedBy=multi-user.target
